# deploy/systemd/disk-cleanup.timer
[Unit]
Description=Limpieza semanal de containerd/buildkit

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
